# settings/config/siem.service
[Unit]
Description=SIEM (Sinatra) — API e dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=siem
Group=siem
WorkingDirectory=/opt/siem
EnvironmentFile=-/opt/siem/.env
# Ajusta RUBY e BUNDLE no PATH do utilizador siem se necessário (rbenv, etc.)
ExecStart=/bin/bash -lc 'cd /opt/siem && bundle exec ruby server.rb'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
